# config/systemd/user/battery-check.service
[Unit]
Description=Battery Level Check

[Service]
Type=oneshot
ExecStart=%h/scripts/batteryAlmostFull.sh
