# scripts/crontab
*/5 * * * * $HOME/scripts/batteryAlmostFull.sh
